# init-files/asaopauloquequeremos.conf
# asaopauloquequeremos ghost

start on startup

script
    cd /var/www/general/asaopauloquequeremos
    npm start --production
end script
