# init-files/seed-ghost.conf
# seed-ghost

start on startup

script
    cd /var/www/general/seed-ghost
    npm start --production
end script
